# supervisor-sshd.ini
[program:sshd]
command=/usr/sbin/sshd -D
prioirty=2

[program:sshd-create-key]
command=/opt/sei-bin/create_ssh_key.sh
autorestart=false
prioirty=1
